# docker/dev/php/supervisord.conf
[program:indexacao_componente_digital]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume indexacao_componente_digital --time-limit=3600 --limit=1000  --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_pessoa]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume indexacao_pessoa --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_modelo]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume indexacao_modelo --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_repositorio]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume indexacao_repositorio --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:indexacao_processo]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume indexacao_processo --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_componente_digital]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume populate_componente_digital --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_pessoa]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume populate_pessoa --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_modelo]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume populate_modelo --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_repositorio]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume populate_repositorio --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:populate_processo]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume populate_processo --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:relatorio_create]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume relatorio_create --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:push_mercure]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume push_mercure --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:push_resource]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume push_resource --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:push_count]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume push_count --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:regras_etiqueta]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume regras_etiqueta --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:integracao_gerar_dossie]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume integracao_gerar_dossie --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:integracao_gerar_dossie_tarefa]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume integracao_gerar_dossie_tarefa --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:process_favorito]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume process_favorito --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:download_processo]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume download_processo --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:envia_ciencia]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume envia_ciencia --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:envia_documento_avulso]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume envia_documento_avulso --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:envia_processo]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume envia_processo --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:envia_componentes_digitais]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume envia_componentes_digitais --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:recebe_recibo_de_tramite]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume recebe_recibo_de_tramite --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:recebe_componentes_digitais]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume recebe_componentes_digitais --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:recebe_tramite]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume recebe_tramite --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:sincroniza_barramento]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume sincroniza_barramento --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:tarefa_barramento]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume tarefa_barramento --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:resposta_oficio_barramento]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume resposta_oficio_barramento --time-limit=3600 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:chat]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume chat --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO

[program:verificacao_virus]
environment=XDEBUG_CONFIG="client_host=192.168.1.5",PHP_IDE_CONFIG="serverName=xdebug-docker",XDEBUG_SESSION="docker",XDEBUG_TRIGGER="1"
command=/usr/bin/php /app/bin/console messenger:consume verificacao_virus --time-limit=3600 --limit=1000 --env=dev
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=AUTO
